# log-clean/first/crontab
00 3 * * * root find /data/logs -maxdepth 1 -type f -mtime +7 -exec rm -rf {} \; && echo "`date` clean OK" >> /var/log/cron.log 2>&1

